# orchestration/netdeneb-auto-update.timer
[Unit]
Description=leggi.es: periodically check for updates

[Timer]
OnCalendar=*:*

[Install]
WantedBy=timers.target
