# orchestration/netdeneb-auto-update.timer
[Unit]
Description=Netdeneb: Every minute, check for Docker image update

[Timer]
OnCalendar=*:*

[Install]
WantedBy=timers.target
