# orchestration/unit-update@.service
[Unit]
Description=Reload systemd and restart '%I'

[Service]
Type=oneshot
ExecStart=systemctl daemon-reload
ExecStart=systemctl restart %I
